# etc/systemd/system/snapraid-btrfs-runner.timer
[Unit]
Description=Run snapraid-btrfs-runner every night

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target

